# scripts/com.us-stock.improve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.us-stock.improve</string>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/frindle/.local/bin/bash-us-stock-improve.sh</string>
		<string>us-stock-improve</string>
		<string>/bin/bash</string>
		<string>/Users/frindle/workspace/synology/us-stock/scripts/run_improvement_cycle.sh</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/frindle/workspace/synology/us-stock/logs/improve.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/frindle/workspace/synology/us-stock/logs/improve.out.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>7</integer>
		<key>Minute</key>
		<integer>0</integer>
		<key>Weekday</key>
		<integer>6</integer>
	</dict>
	<key>WakeForScheduledEvents</key>
	<true/>
</dict>
</plist>
