# scripts/com.mpva.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mpva.monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-i</string>
        <string>/Users/frindle/workspace/synology/us-stock/.venv/bin/python3</string>
        <string>/Users/frindle/workspace/synology/us-stock/docs/_ideas/vending-machine/monitor.py</string>
    </array>

    <!-- 3시간마다 실행: 0, 3, 6, 9, 12, 15, 18, 21시 (주말 포함) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/frindle/workspace/synology/us-stock/logs/mpva_launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/frindle/workspace/synology/us-stock/logs/mpva_launchd.err.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/frindle</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>TELEGRAM_BOT_TOKEN</key>
        <string>8710179946:AAGiYmKYINInCk8bQJ1xYBaFlmPI73KqycQ</string>
        <key>TELEGRAM_CHAT_ID</key>
        <string>183678616</string>
    </dict>
</dict>
</plist>
